# systemd/dms-wallpaper-watcher.service
[Unit]
Description=DMS Wallpaper Watcher - Updates awww and matugen when DMS wallpaper changes
After=dms.service
BindsTo=dms.service

[Service]
Type=simple
ExecStart=/home/higorprado/.local/bin/dms-wallpaper-watcher
Restart=always
RestartSec=5

[Install]
WantedBy=dms.service
